# containers/lifeos-llama-server/lifeos-llama-server.container
# lifeos-llama-server.container — Qwen3.5 Vulkan inference Quadlet
# Default LIFEOS_AI_MODEL is Qwen3.5-4B (per decision_default_model_4b);
# 9B is opt-in via the runtime-profile env file. The Quadlet itself is
# model-agnostic.
#
# Phase 4 of the architecture pivot. NOT auto-installed; BLOCKED on
# nvidia-container-toolkit availability in the bootc image (see
# project_pending_nvidia_container_toolkit memory).

[Unit]
Description=LifeOS LLM Server (Qwen3.5 chat, Vulkan GPU, containerized)
Documentation=https://github.com/ggml-org/llama.cpp
# lifeosd reads from :8082 — keep llama-server up before lifeosd. Use the
# Quadlet-generated unit name; ordering vs the legacy `llama-server.service`
# (host) is silently dropped because that unit is masked post-Phase-4.
After=network-online.target lifeos-image-guardian.service lifeos-net-network.service
# Hard bridge dependency — see lifeos-tts.container for rationale.
Requires=lifeos-net-network.service
Before=lifeos-lifeosd.service
# Note: NO After=/Wants=nvidia-cdi-refresh.service.
# nvidia-cdi-refresh.service ships from nvidia-container-toolkit with
# `After=multi-user.target`. Our Quadlet has `WantedBy=multi-user.target`,
# so adding After= or Wants= against it produces an ordering cycle that
# systemd breaks by DELETING our service from the boot transaction —
# observed live in production: Phase-4 deploy left llama-server failing on
# every boot until manually restarted, with the journal showing
# "Found ordering cycle ... Job lifeos-llama-server.service/start deleted
# to break ordering cycle". The CDI refresh runs independently when the
# .path unit triggers (modules.dep / nvidia-ctk install changes). The
# ExecStartPre below waits for the CDI spec to actually exist before
# podman tries to inject the device, which races cleanly with whenever
# nvidia-cdi-refresh decides to write it.
Wants=network-online.target

[Container]
Image=ghcr.io/hectormr206/lifeos-llama-server:stable
ContainerName=lifeos-llama-server
HostName=lifeos-llama-server

# Phase 8: bridge network. The Exec= below now binds the server on
# `--host 0.0.0.0` inside the container (container netns) so it listens on
# the bridge IP; PublishPort= forwards host loopback :8082 to that bridge IP.
# lifeosd (still on Network=host) keeps using http://127.0.0.1:8082 unchanged.
Network=lifeos-net.network
PublishPort=127.0.0.1:8082:8082

# === GPU passthrough via CDI ===
#
# AddDevice=nvidia.com/gpu=all uses Container Device Interface. The host
# nvidia-container-toolkit generates /etc/cdi/nvidia.yaml that describes how
# to inject /dev/nvidia*, libraries, and Vulkan ICDs into the container.
# This replaces the legacy `--gpus all` Docker flag with a runtime-agnostic
# spec.
#
# Requirements on host:
#   - nvidia-container-toolkit installed (BLOCKED in current image — see
#     memory project_pending_nvidia_container_toolkit)
#   - nvidia-cdi-refresh.service ran successfully (ships with the toolkit)
#   - GPU is Turing+ if using nvidia-open driver (we use closed via
#     lifeos-nvidia-drivers — works for Pascal+)
AddDevice=nvidia.com/gpu=all

# GPU access mechanism is CDI device injection (the AddDevice= line
# above), NOT group membership. Round-2 audit caught the prior version
# claiming GroupAdd=keep-groups was the fix — that's wrong: keep-groups
# propagates the rootless caller's supplementary groups, but for a
# rootful Quadlet the caller is systemd PID 1 running as root, whose
# supplementary group set is empty.
#
# CDI's spec (/etc/cdi/nvidia.yaml, generated by nvidia-ctk) injects:
#   - /dev/nvidia* device files (with mode 666 by default — kernel
#     skips the group check)
#   - cgroup device permission rules that grant rwm at the cgroup layer
#   - libnvidia-* shared libraries inside the container's lib paths
#   - NVIDIA Vulkan ICD JSON
# None of those require the container user to be in a `nvidia` group.
#
# If a future driver tightens /dev/nvidia* permissions to mode 660 with
# group=nvidia, the right fix is `GroupAdd=video,render` (which podman
# resolves against host /etc/group at start) or explicit numeric GIDs —
# NOT keep-groups.

# === Models bind mount — the largest single state item in LifeOS ===
#
# /var/lib/lifeos/models holds:
#   Qwen3.5-4B-Q4_K_M.gguf      ~2.5 GB ← default model (per
#                                          decision_default_model_4b)
#   Qwen3.5-4B-mmproj-F16.gguf  ~640 MB ← default multimodal projection
#   + opt-in larger 9B (5.3 GB + 876 MB mmproj) for hosts with 8GB+ VRAM
#   + smaller 2B/0.8B fallbacks for low-spec hosts
#   + nomic-embed-text-v1.5.f16.gguf  (used by lifeos-llama-embeddings)
#
# Always bind-mount RO. Models NEVER live inside the container image — that
# would inflate to ~14 GB and force a rebuild on model bumps.
#
# ':z' (lowercase, SHARED label) — same dir is also mounted by lifeosd
# and lifeos-llama-embeddings. ':Z' would assign a per-container private
# MCS category and lock the others out. ':z' shares cleanly.
Volume=/var/lib/lifeos/models:/var/lib/lifeos/models:ro,z

# Runtime profile env (LIFEOS_AI_MODEL, LIFEOS_AI_CTX_SIZE, LIFEOS_AI_THREADS,
# LIFEOS_AI_GPU_LAYERS, LIFEOS_AI_PARALLEL, etc.) lives on host so the
# benchmarker can update it without rebuilding the container.
#
# NO Volume= bind for these env files. The [Service] EnvironmentFile=
# directives below load them on the HOST side (systemd reads, exports
# vars to the podman process, podman injects as ENV into the container).
# A separate bind mount would be redundant — llama-server reads vars
# from its environment, not from a file path.

# === The exec line — same args the host service uses ===
#
# Reads model + tuning from the runtime-profile env file. The benchmarker
# (defer benchmark, PR #65) writes that file with hardware-tuned values.
# Game-guard env file overrides for low-power profiles when gaming detected.
Exec=/usr/sbin/llama-server \
    --model /var/lib/lifeos/models/${LIFEOS_AI_MODEL} \
    --mmproj /var/lib/lifeos/models/${LIFEOS_AI_MMPROJ} \
    --alias lifeos \
    --host 0.0.0.0 \
    --port 8082 \
    --ctx-size ${LIFEOS_AI_CTX_SIZE} \
    --threads ${LIFEOS_AI_THREADS} \
    --n-gpu-layers ${LIFEOS_AI_GPU_LAYERS} \
    --mlock \
    --parallel ${LIFEOS_AI_PARALLEL} \
    --batch-size ${LIFEOS_AI_BATCH_SIZE} \
    --ubatch-size ${LIFEOS_AI_UBATCH_SIZE} \
    --n-predict 2048 \
    --flash-attn auto \
    --cache-type-k q8_0 \
    --cache-type-v q8_0 \
    -kvu -cram -1 -sm none -mg 0 \
    --temp 0.6 --top-p 0.95 --top-k 20 --min-p 0.0 \
    --presence-penalty 0.0 --repeat-penalty 1.0 \
    --reasoning off \
    --jinja

[Service]
# Defense-in-depth Capa 4: ensure image is present locally before starting.
# Quadlet requires ExecStartPre in [Service], not [Container].
ExecStartPre=/usr/local/bin/lifeos-ensure-images

# Wait up to 60 s for nvidia-cdi-refresh.service to publish the CDI spec.
# We can't `After=` it (ordering cycle — see [Unit] comments), so we poll
# the canonical paths instead. NVIDIA's unit writes /var/run/cdi/nvidia.yaml
# (NVIDIA_CTK_CDI_OUTPUT_FILE_PATH default); /etc/cdi/ is a manual path
# operators sometimes populate. Either lets podman's CDI lookup succeed
# for `AddDevice=nvidia.com/gpu=all`. Hyphen-prefix + bash makes the step
# best-effort: if 60 s pass without a spec we still try to start (Restart=
# always handles a hard failure).
ExecStartPre=-/bin/bash -c 'i=0; while [ $i -lt 60 ]; do for f in /run/cdi/nvidia.yaml /var/run/cdi/nvidia.yaml /etc/cdi/nvidia.yaml; do [ -s "$f" ] && exit 0; done; sleep 1; i=$((i+1)); done; exit 0'

# CRITICAL: these paths are HOST paths (systemd evaluates EnvironmentFile=
# on the host before the container is created). The benchmarker writes
# to /var/lib/lifeos/llama-server-runtime-profile.env on the host —
# previous version pointed at /etc/llama-server-runtime-profile.env which
# does NOT exist on the host, so all ${LIFEOS_AI_*} vars came up empty
# and the Exec= line above ran with --model "/var/lib/lifeos/models/"
# (no filename) → llama-server fails to start.
EnvironmentFile=-/var/lib/lifeos/llama-server-runtime-profile.env
EnvironmentFile=-/var/lib/lifeos/llama-server-game-guard.env

# Capa 3 — auto-restart with a longer backoff than embeddings/TTS because
# Qwen3.5 warmup (model load + GPU upload + Vulkan shader cache prime)
# takes ~30-60s. Avoid restart-thrash on transient OOM.
Restart=always
RestartSec=15s
TimeoutStartSec=180s

# Resource limits. The model itself fills GPU VRAM (5.3 GB Q4); host RAM is
# for the kv-cache, server overhead, and Vulkan staging. 8G cap is generous
# but realistic for the 131K context profile.
MemoryMax=8G
CPUQuota=400%

# DO NOT add `NoNewPrivileges=true`. Same Fedora 44 + podman 5.x
# `--cgroups=split` interaction that broke TTS in PR #77 — NNP at the
# [Service] layer blocks /run/libpod/alive.lck acquisition. The other
# four LifeOS Quadlets have this comment too; this one was the outlier
# until JD round 1 caught it.

[Install]
WantedBy=multi-user.target default.target
