# image/files/etc/systemd/system/spice-vdagentd.service.d/10-lifeos-vm-only.conf
# LifeOS — Service Minimization (CIS Benchmark)
# spice-vdagent only needed inside VMs (clipboard/resize).
# ConditionVirtualization ensures it never starts on bare metal.
[Unit]
ConditionVirtualization=vm
