# image/files/etc/systemd/system/whisper-stt.service.d/20-lifeos-background.conf
# Pin whisper-stt to the LifeOS background slice (E-cores, low priority).
# STT transcription is not latency-sensitive — running on E-cores saves
# P-core power budget for interactive tasks.
[Service]
Slice=lifeos-background.slice
Nice=15
IOSchedulingClass=idle
