# image/files/etc/udev/rules.d/85-lifeos-power-epp.rules
# LifeOS — Dynamic Energy Performance Preference switching.
#
# On AC: balance_performance (responsive, uses P-cores freely)
# On battery: balance_power (prefers E-cores, conservative boosting)
#
# This is the single biggest software lever for battery life on Intel
# hybrid CPUs. The EPP hint tells intel_pstate how aggressively to
# boost — balance_power alone can save 20-30% power on battery.
#
# Works on any Intel CPU with HWP (Hardware P-states) support.
# AMD equivalent (amd-pstate EPP) uses the same sysfs path.
#
# NOTE: udev's RUN= parser rejects `$` as anything but a udev substitution,
# so we cannot inline a shell loop with `$f`. The real logic lives in
# /usr/local/bin/lifeos-epp-switch.sh — udev just invokes it.

# AC plugged in → responsive
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", ATTR{online}=="1", RUN+="/usr/local/bin/lifeos-epp-switch.sh balance_performance"

# AC unplugged → efficient
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", ATTR{online}=="0", RUN+="/usr/local/bin/lifeos-epp-switch.sh balance_power"
