# image/files/etc/udev/rules.d/90-lifeos-thermal.rules
# Allow lifeos user to control CPU performance scaling for thermal management.
# This enables the daemon's reactive thermal manager to adjust caps dynamically
# without requiring sudo or NoNewPrivileges exceptions.
#
# Note: intel_pstate and cpufreq policy* nodes are virtual sysfs files, not
# udev devices themselves. We hook on per-CPU "cpu" subsystem add/change
# events (which DO fire at boot and on hotplug for every cpuN) and chmod the
# runtime paths. The `test -f` guards make this safe when a driver is absent.
#
# Ownership: root:lifeos 0660 — least-privilege (ref: reference_sudo_policy).
# The daemon runs as uid 1000 which is the primary member of group `lifeos`.
#
# NOTE: udev's RUN= parser rejects `$` as anything but a udev substitution,
# so we cannot inline a shell loop with `$f`. The real logic lives in
# /usr/local/bin/lifeos-thermal-grant.sh — udev just invokes it.

# Per-CPU event — fires for every cpu0..cpuN at boot and on hotplug.
SUBSYSTEM=="cpu", ACTION=="add|change", RUN+="/usr/local/bin/lifeos-thermal-grant.sh"
