# image/files/etc/systemd/system/lifeos-sentinel.service.d/10-dev-mode-override.conf
# LifeOS dev-mode: override ExecStart to a writable copy of the sentinel
# script. /usr/local/bin is part of /usr which is read-only on bootc, so
# iterating on the sentinel during dev requires staging the script in
# /var/lib/lifeos/bin/ and pointing the unit there until the next image
# rebuild bakes the fix into /usr/local/bin/.
#
# The script is executed via /bin/bash because SELinux labels scripts
# under /var/lib/ with var_lib_t, which systemd refuses to exec directly
# (EACCES on load). Running bash (labelled bin_t) as the interpreter
# sidesteps the label check without weakening the targeted policy.
#
# Release images ship the original script under /usr/local/bin and must
# NOT install this dropin. The image build verification block in
# image/Containerfile is responsible for keeping release images clean.
[Service]
ExecStart=
ExecStart=/usr/bin/bash -c 'if [ -f /var/lib/lifeos/bin/lifeos-sentinel.sh ]; then exec /usr/bin/bash /var/lib/lifeos/bin/lifeos-sentinel.sh; else exec /usr/bin/bash /usr/local/bin/lifeos-sentinel.sh; fi'
