# image/files/etc/systemd/system/systemd-remount-fs.service.d/10-skip-on-overlay.conf
[Service]
ExecCondition=
ExecCondition=/bin/sh -c '[ "$(findmnt -n -o FSTYPE /)" != "overlay" ]'
