# containers/lifeos-simplex-bridge/lifeos-simplex-bridge.container
# lifeos-simplex-bridge.container — SimpleX CLI bot Quadlet

[Unit]
Description=LifeOS SimpleX Chat Bot (Axi private messaging — containerized)
Documentation=https://github.com/simplex-chat/simplex-chat
After=network-online.target lifeos-image-guardian.service
Wants=network-online.target
# lifeosd connects to :5226 — keep this up before lifeosd.
Before=lifeos-lifeosd.service

[Container]
Image=ghcr.io/hectormr206/lifeos-simplex-bridge:stable
ContainerName=lifeos-simplex-bridge
HostName=lifeos-simplex-bridge

# Phase 8: bridge network. simplex-chat binds its WebSocket on :5226 inside
# the container netns; PublishPort= forwards host loopback :5226 to the
# bridge IP. lifeosd (Network=host) keeps its existing
# ws://127.0.0.1:5226 URL. SimpleX protocol's outbound TCP/443 to SMP relays
# works equally on bridge — podman SNAT handles it.
Network=lifeos-net.network
PublishPort=127.0.0.1:5226:5226

# === SimpleX state — survives container recreate ===
#
# The CLI creates several SQLite files under the bind-mounted prefix:
#   bot_chat.db, bot_files.db, bot_agent.db, etc.
# Plus the message vault (encrypted local store) and the contact list.
# All of this is the single most important state piece for the bridge —
# losing it means re-pairing every contact (Hector's phone, etc.).
#
# ':z' (lowercase, shared MCS label). The legacy host
# `simplex-chat.service` reads the same path during rollback flow;
# ':Z' (private) would lock the legacy service out post-rollback.
Volume=/var/lib/lifeos/simplex:/var/lib/lifeos/simplex:z

# CA bundle for TLS to SMP relays.
# ':z' shared label — touching the system-wide CA trust dir with ':Z'
# would relabel it and break TLS for every other host service.
# NO `:z` on /etc/pki/ca-trust/extracted — relabeling host CA bundle to
# container_file_t breaks TLS for everything else on the host (curl,
# openssl, etc. expect cert_t / system_conf_t). Plain `:ro` is correct;
# Fedora ships /etc/pki readable by container_t by default.
Volume=/etc/pki/ca-trust/extracted:/etc/pki/ca-trust/extracted:ro

[Service]
# Defense-in-depth Capa 4: ensure image is present locally before starting.
# Quadlet requires ExecStartPre in [Service], not [Container].
ExecStartPre=/usr/local/bin/lifeos-ensure-images

# Provision /var/lib/lifeos/simplex with the right perms (legacy unit's
# StateDirectory= no longer runs since we removed it from
# multi-user.target.wants). Idempotent — install -d only creates if absent.
ExecStartPre=/usr/bin/install -d -m 0700 /var/lib/lifeos/simplex

# Bootstrap the SimpleX bot profile (Axi). The script runs simplex-chat
# once with `-e /quit` to satisfy the interactive "display name:" prompt
# and seed the bot DB. Idempotent — exits 0 if the marker file already
# exists. Same script the legacy simplex-chat.service used.
ExecStartPre=/usr/local/bin/lifeos-simplex-setup.sh

# Capa 3 — SimpleX is touchy about graceful shutdown (DB locks). RestartSec
# matches the host service's 10s, giving SimpleX time to flush WAL on SIGTERM.
Restart=always
RestartSec=10s
TimeoutStartSec=30s

# Resource limits. SimpleX is light: ~50-100 MB RAM, minimal CPU.
MemoryMax=512M
CPUQuota=100%

# DO NOT add `NoNewPrivileges=true` here — same Fedora podman --cgroups=split
# interaction that broke TTS. Container-side hardening covers what NNP would
# at the systemd layer.

[Install]
WantedBy=multi-user.target default.target
