# files/usr/lib/systemd/system/lifeos-health.timer
[Unit]
Description=LifeOS Health Check Timer
Documentation=man:lifeosd(8)

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
