# files/usr/lib/systemd/system/lifeos-health.service
[Unit]
Description=LifeOS Health Check Service
Documentation=man:lifeosd(8)

[Service]
Type=oneshot
ExecStart=/usr/bin/life osd health-check
StandardOutput=journal
StandardError=journal
