# files/usr/lib/systemd/system/lifeosd.service
[Unit]
Description=LifeOS System Daemon
Documentation=man:lifeosd(8)
After=network.target
Wants=network.target

[Service]
Type=notify
ExecStart=/usr/bin/lifeosd
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lifeosd

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/lifeos /var/log/lifeos
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
