# containers/lifeos-llama-embeddings/lifeos-llama-embeddings.container
# lifeos-llama-embeddings.container — Quadlet for nomic-embed-text server
#
# Phase 2 of the architecture pivot. NOT auto-installed yet — see
# containers/lifeos-llama-embeddings/README.md for the manual deploy recipe.

[Unit]
Description=LifeOS Semantic Embeddings Server (containerized)
Documentation=https://huggingface.co/nomic-ai/nomic-embed-text-v1.5-GGUF
After=network-online.target lifeos-image-guardian.service lifeos-net-network.service
Wants=network-online.target
# Hard bridge dependency — see lifeos-tts.container for rationale.
Requires=lifeos-net-network.service
# lifeosd does memory recall via embeddings — keep this up before lifeosd.
# Use the Quadlet-generated unit name; ordering vs the legacy
# 'lifeosd.service' is silently dropped because that unit doesn't exist
# after Phase 3 flips.
Before=lifeos-lifeosd.service

[Container]
Image=ghcr.io/hectormr206/lifeos-llama-embeddings:stable
# Phase 6 — auto-update: podman-auto-update.timer (daily) compares the
# registry digest for the pinned tag and rolling-restarts the service if
# it changed.
AutoUpdate=registry
ContainerName=lifeos-llama-embeddings
HostName=lifeos-llama-embeddings

# Phase 8: bridge network. Server binds 0.0.0.0:8083 inside the container
# netns; PublishPort= forwards host loopback :8083 → bridge IP. lifeosd
# (Network=host) keeps reaching it via http://127.0.0.1:8083 unchanged.
Network=lifeos-net.network
PublishPort=127.0.0.1:8083:8083

# Models bind mount — read-only. The model file (currently
# nomic-embed-text-v1.5.f16.gguf, ~262 MB) lives on host /var/lib/lifeos/models/
# and survives container recreation. NEVER bake the model into the image —
# it'd inflate the image to ~400+ MB and forces a rebuild on model bumps.
#
# ':z' (lowercase, SHARED label). /var/lib/lifeos/models is read by
# lifeos-llama-embeddings AND lifeos-llama-server AND lifeosd. ':Z'
# (uppercase, private) would assign a per-container MCS category and
# lock out the others. ':z' shares the label cleanly across consumers.
Volume=/var/lib/lifeos/models:/models:ro,z

# Env file from host so model name (LIFEOS_EMBED_MODEL) is configurable
# without rebuilding the container.
# ':z' shared label — /etc/lifeos is read by multiple containers + the
# legacy host service (rollback path). Same reasoning as models bind.
Volume=/etc/lifeos/llama-embeddings.env:/etc/lifeos/llama-embeddings.env:ro,z

# Override the image's CMD to point at the model from the bind mount. The
# model name is read from the env file at start time — keeps the Quadlet
# generic across model bumps.
Exec=/usr/sbin/llama-server \
    --model /models/${LIFEOS_EMBED_MODEL} \
    --alias lifeos-embeddings \
    --host 0.0.0.0 \
    --port 8083 \
    --ctx-size 2048 \
    --threads 2 \
    --batch-size 256 \
    --ubatch-size 64 \
    --embeddings \
    --pooling mean

[Service]
# Defense-in-depth Capa 4: ensure image is present locally before starting.
# Quadlet requires ExecStartPre in [Service], not [Container].
ExecStartPre=/usr/local/bin/lifeos-ensure-images

# Bootstrap the GGUF model + write LIFEOS_EMBED_MODEL into the env file.
# Leading-dash (`-`) makes systemd ignore non-zero exits from this step.
# The legacy host service used `ExecCondition=`, which has the semantic
# that exit 1-254 is a CLEAN SKIP. Quadlet only allows ExecStartPre, so
# we replicate the skip-on-non-zero behavior with the dash prefix —
# essential for the model-missing + LIFEOS_AI_AUTO_MANAGE_MODELS=false
# path (script exits 1; without the dash the unit would crash-loop).
# Script contract: 0 = ready, 1-254 = clean skip, 255 = hard fail.
ExecStartPre=-/usr/local/bin/lifeos-embeddings-setup.sh

EnvironmentFile=-/etc/lifeos/llama-embeddings.env

# Capa 3 of defense-in-depth — auto-restart with sane backoff. Embeddings
# warm-up is fast (~5s), so RestartSec is shorter than TTS (which warms ~15s).
Restart=always
RestartSec=5s
TimeoutStartSec=60s

# Resource limits — slightly tighter than the host service (1G) because
# container overhead is ~50-100 MB extra so 1280M caps the total at the same
# ~1.3 GB envelope.
MemoryMax=1280M
CPUQuota=200%

# DO NOT add `NoNewPrivileges=true` here — same Fedora podman --cgroups=split
# interaction that broke TTS on first reboot. The Quadlet generator already
# applies container-side isolation; NNP at the [Service] layer is redundant
# and breaks runtime init lock acquisition.

[Install]
WantedBy=multi-user.target default.target
