# image/files/etc/systemd/user/lifeosd.service.d/00-tts.conf
[Service]
# URL of the local Kokoro TTS HTTP server (see lifeos-tts-server.service).
# Without this, sensory_pipeline.rs::synthesize_tts falls back to espeak-ng and
# every voice sounds like a 1980s robot regardless of the user's selection.
Environment=LIFEOS_TTS_SERVER_URL=http://127.0.0.1:8084

# Default voice used when no per-user preference is set. Must match a voice
# returned by GET /voices on the TTS server.
Environment=LIFEOS_TTS_DEFAULT_VOICE=ef_dora
