# image/files/etc/systemd/user/lifeosd.service.d/10-dev-mode.conf
# LifeOS Developer Mode — lifeosd environment override
#
# ONLY installed when LIFEOS_BUILD_MODE=dev at image build time.
# Release images do NOT contain this file.
#
# Sets RUST_LOG=debug by default so the AI assistant has verbose output
# to work with when diagnosing runtime issues. Individual crates can
# still be filtered at runtime via the usual env_logger syntax.
#
# To disable without rebuilding:
#     systemctl --user edit lifeosd
# and add:
#     [Service]
#     Environment=RUST_LOG=info

[Service]
Environment=RUST_LOG=debug
Environment=LIFEOS_DEV_MODE=1
