# image/files/etc/systemd/system/conduit.service
[Unit]
Description=Conduit Matrix Homeserver (LifeOS)
Documentation=https://conduit.rs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=CONDUIT_CONFIG=/etc/lifeos/conduit.toml
ExecStartPre=/usr/local/bin/lifeos-matrix-setup.sh
ExecStart=/usr/sbin/conduit
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes

# Conduit needs to write its database and Axi credentials
ReadWritePaths=/var/lib/lifeos/conduit /etc/lifeos

[Install]
WantedBy=multi-user.target
