# image/files/etc/systemd/system/lifeos-update-check.timer
[Unit]
Description=LifeOS Daily Update Check Timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
