# image/files/etc/systemd/system/lifeos-update-check.service
[Unit]
Description=LifeOS Daily Update Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lifeos-update-check.sh
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/var/lib/lifeos
