# image/files/etc/systemd/system/lifeos-sentinel.service
[Unit]
Description=LifeOS Sentinel — Independent watchdog for lifeosd
After=network-online.target lifeosd.service
Wants=lifeosd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/lifeos-sentinel.sh
Restart=always
RestartSec=10

# Security: minimal privileges
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log/lifeos
ProtectHome=true

[Install]
WantedBy=multi-user.target
