# image/files/etc/systemd/system/lifeos-btrfs-snapshot.service
[Unit]
Description=LifeOS Btrfs Snapshot
Documentation=man:btrfs-subvolume(8)
ConditionPathExists=/usr/bin/btrfs

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lifeos-btrfs-snapshot.sh scheduled
User=root
Group=root
