# image/files/etc/systemd/system/lifeos-smart-charge.timer
[Unit]
Description=LifeOS Smart Charging Schedule
Documentation=https://github.com/hectormr206/lifeos

[Timer]
# Run every hour to adjust charge threshold based on time of day
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
