# image/files/etc/systemd/system/lifeos-tts-server.service
[Unit]
Description=LifeOS Kokoro TTS Server (HTTP, 127.0.0.1:8084)
Documentation=https://github.com/hexgrad/kokoro
After=network.target
Before=lifeosd.service

[Service]
Type=simple

# Runtime environment
EnvironmentFile=/etc/lifeos/tts-server.env

# Server binary
ExecStart=/opt/lifeos/kokoro-env/bin/python3 /usr/local/bin/lifeos-tts-server.py

# Restart policy
# `always` because the memory watchdog inside the script exits cleanly (status 0)
# after SIGTERM'ing itself; `on-failure` would leave the server dead. StartLimit
# prevents infinite restart loops if the process genuinely cannot stay within the limit.
Restart=always
RestartSec=5
TimeoutStartSec=60
StartLimitBurst=3
StartLimitIntervalSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lifeos-tts-server

# Dynamic user for isolation (no persistent UID needed)
DynamicUser=yes
RuntimeDirectory=lifeos-tts

# Memory and CPU limits
# 1536M leaves ~140MB cushion above the in-process watchdog (1400MB) so the kernel
# OOM killer never fires before the script's graceful SIGTERM path.
MemoryMax=1536M
CPUQuota=200%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native
LockPersonality=true
CapabilityBoundingSet=

# Network isolation — loopback only
IPAddressAllow=127.0.0.0/8 ::1/128
IPAddressDeny=any
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
