# image/files/etc/systemd/system/lifeos-state-migrate.service
[Unit]
Description=LifeOS State Migration — ~/.local/share/lifeos to /var/lib/lifeos
Documentation=man:lifeos-state-migrate(8)
# No ConditionPathExists= for the source: the script resolves UID 1000's
# home dynamically via getent, so a hard-coded /home/lifeos path would
# wrongly skip the unit on hosts where the LifeOS user has a different name.
# The script itself exits cleanly when there's nothing to migrate.
# Run BEFORE the lifeos-lifeosd Quadlet so the data is in place when the
# container opens memory.db. systemd silently drops Before= against units
# that don't exist (e.g. on a fresh install where lifeos-lifeosd.service
# never gets generated yet), so this is safe in all configurations.
Before=lifeos-lifeosd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/lifeos-state-migrate.sh
TimeoutStartSec=30m
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lifeos-state-migrate

[Install]
WantedBy=multi-user.target
