# image/files/etc/systemd/system/lifeos-ecore-pin.service
[Unit]
Description=LifeOS E-Core Pin — restrict background slice to efficiency cores
# NOTE: we MUST NOT order After=multi-user.target while also Before= units
# that are WantedBy=multi-user.target — that creates an ordering cycle that
# systemd breaks by deleting start jobs for llama-embeddings/whisper-stt.
# Instead we order after basic.target (cgroups ready) and before the
# background workloads themselves; multi-user.target still pulls us in
# through the [Install] section.
After=basic.target
Before=whisper-stt.service llama-embeddings.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Detect E-cores via the standalone script and apply cpuset to the
# background slice. On non-hybrid CPUs this sets all CPUs (no restriction).
ExecStart=/bin/bash -c '\
  SCRIPT=/usr/local/bin/lifeos-detect-ecores.sh; \
  [ -x "$$SCRIPT" ] || SCRIPT=/etc/lifeos/lifeos-detect-ecores.sh; \
  if [ -x "$$SCRIPT" ]; then \
    ECORES=$$($$SCRIPT); \
  else \
    ECORES="0-$$(($$(nproc)-1))"; \
    echo "[ecore-pin] No detect script found — using all CPUs"; \
  fi; \
  echo "[ecore-pin] E-cores: $${ECORES}"; \
  systemctl set-property lifeos-background.slice AllowedCPUs="$${ECORES}" 2>/dev/null || \
    echo "[ecore-pin] Could not set AllowedCPUs (cgroup v2 required)"'

[Install]
WantedBy=multi-user.target
