# image/files/etc/systemd/system/llama-embeddings.service
[Unit]
Description=LifeOS Semantic Embeddings Server (nomic-embed-text via llama.cpp)
Documentation=https://huggingface.co/nomic-ai/nomic-embed-text-v1.5-GGUF
After=network-online.target systemd-user-sessions.service
Wants=network-online.target

[Service]
Type=simple
# Inherit LIFEOS_AI_AUTO_MANAGE_MODELS from the shared AI config.
EnvironmentFile=-/etc/lifeos/llama-server.env
EnvironmentFile=-/etc/lifeos/llama-embeddings.env
ExecStartPre=/usr/local/bin/lifeos-embeddings-setup.sh
ExecStart=/usr/sbin/llama-server \
    --model /var/lib/lifeos/models/${LIFEOS_EMBED_MODEL} \
    --alias lifeos-embeddings \
    --host 127.0.0.1 \
    --port 8083 \
    --ctx-size 2048 \
    --threads 2 \
    --batch-size 256 \
    --ubatch-size 64 \
    --embeddings \
    --pooling mean
Restart=always
RestartSec=15
TimeoutStartSec=60
StartLimitInterval=300
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llama-embeddings

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/lifeos /etc/lifeos
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=32768
# nomic-embed-text-v1.5 is ~84MB on disk; cap to 1G to leave headroom for batching.
MemoryMax=1G

[Install]
WantedBy=multi-user.target
