# image/files/etc/systemd/system/llama-server.service
[Unit]
Description=LifeOS AI Inference Server (llama.cpp)
Documentation=https://github.com/ggml-org/llama.cpp
After=network.target
Wants=network.target
ConditionPathExists=/etc/lifeos/llama-server.env

[Service]
Type=simple
EnvironmentFile=-/etc/lifeos/llama-server.env
ExecStartPre=/usr/local/bin/lifeos-ai-setup.sh
ExecStart=/usr/sbin/llama-server \
    --model /var/lib/lifeos/models/${LIFEOS_AI_MODEL} \
    --mmproj /var/lib/lifeos/models/${LIFEOS_AI_MMPROJ} \
    --alias ${LIFEOS_AI_ALIAS} \
    --host ${LIFEOS_AI_HOST} \
    --port ${LIFEOS_AI_PORT} \
    --ctx-size ${LIFEOS_AI_CTX_SIZE} \
    --threads ${LIFEOS_AI_THREADS} \
    --n-gpu-layers ${LIFEOS_AI_GPU_LAYERS} \
    --jinja
Restart=on-failure
RestartSec=10
TimeoutStartSec=120
StartLimitInterval=300
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llama-server

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateUsers=true
ReadWritePaths=/var/lib/lifeos /etc/lifeos
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native
SystemCallFilter=@system-service

# Resource limits
LimitNOFILE=65536
MemoryMax=16G

# GPU device access (conditional)
DeviceAllow=/dev/nvidia* rw
DeviceAllow=/dev/kfd rw
DeviceAllow=/dev/dri/* rw
DeviceAllow=/dev/dxg rw

[Install]
WantedBy=multi-user.target
