# image/files/etc/systemd/system/ollama.socket
[Unit]
Description=Ollama AI Service Socket
Documentation=https://github.com/ollama/ollama

[Socket]
ListenStream=0.0.0.0:11434
Accept=false
NoDelay=true
DeferAcceptSec=1

[Install]
WantedBy=sockets.target
