# image/files/etc/systemd/system/whisper-stt.service
[Unit]
Description=LifeOS STT Runtime (whisper.cpp)
Documentation=https://github.com/ggml-org/whisper.cpp
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/lifeos-stt-setup.sh
ExecStop=/usr/bin/true
TimeoutStartSec=600

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whisper-stt

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/lifeos
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
