# image/files/etc/systemd/system/libvirtd.service.d/10-lifeos-socket-only.conf
# LifeOS — Service Minimization (CIS Benchmark)
# libvirtd should only start via socket activation, not at boot.
# This saves ~30MB RAM when no VMs are running.
[Unit]
ConditionPathExists=/var/lib/libvirt
