# image/files/etc/systemd/system/nvidia-persistenced.service.d/10-lifeos-conditions.conf
[Unit]
# Avoid failed state on systems where NVIDIA userspace exists but the proprietary
# kernel module is not active (e.g. nouveau fallback or unsigned modules).
ConditionPathExists=/dev/nvidiactl
