# containers/lifeos-tts/lifeos-tts.container
# lifeos-tts.container — Quadlet for Kokoro TTS server (self-contained image)
#
# Phase 6 of the architecture pivot — the side image now carries the full
# Kokoro venv + voice models + voices-manifest.json + the server script.
# The bootc image no longer ships /opt/lifeos/kokoro-env or
# /usr/local/bin/lifeos-tts-server.py — both moved INTO the container
# image, removing ~850 MB from bootc.
#
# PRD: docs/strategy/prd-architecture-pivot-lean-bootc-quadlet.md

[Unit]
Description=LifeOS Kokoro TTS Server (containerized)
Documentation=https://github.com/hexgrad/kokoro
After=network-online.target lifeos-image-guardian.service lifeos-net-network.service
Wants=network-online.target
# Hard dependency on the bridge: without it the container starts but
# can't attach to lifeos-net, PublishPort= forwards die silently, and
# the daemon's loopback connections are refused. Requires= pulls the
# network in if it's not already running.
Requires=lifeos-net-network.service
# Order before the daemon Quadlet so TTS is up by the time the daemon's
# sensory_pipeline first probes /tts. Phase 3 made lifeosd itself a
# Quadlet (`lifeos-lifeosd.service`), so this Before= now resolves
# against a real generated unit instead of a no-op.
Before=lifeos-lifeosd.service

[Container]
Image=ghcr.io/hectormr206/lifeos-tts:stable
ContainerName=lifeos-tts
HostName=lifeos-tts

# Phase 8: container moves off Network=host onto the dedicated podman bridge
# `lifeos-net`. Publish=127.0.0.1:8084:8084 keeps the daemon's existing
# loopback URL (http://127.0.0.1:8084) working unchanged — the daemon is
# still on Network=host until its hardcoded URLs (sensory_pipeline.rs,
# llm_router.rs, ai.rs) move behind env vars that can resolve container DNS
# (Phase 8b daemon refactor — deferred). Net effect today: the TTS process
# no longer shares the host network namespace, which is the defense-in-depth
# win Phase 8 was about. End-to-end behaviour unchanged.
Network=lifeos-net.network
PublishPort=127.0.0.1:8084:8084

# Only surface left: the env file. Phase 6 makes everything else
# (venv, voice models, manifest, server script) live INSIDE the
# container image.
#
# We use Quadlet's [Container] EnvironmentFile= which translates to
# `podman --env-file=/etc/lifeos/tts-server.env`, propagating the
# file's vars into the container process env. The script reads via
# os.environ, never opens the path inside the container — so a
# Volume= bind of the same file would be redundant AND would trigger
# an SELinux `:z` relabel of /etc/lifeos/tts-server.env to
# container_file_t (host services that read the same file would
# then get an unexpected MCS label). Single source of truth: this
# EnvironmentFile= line.
EnvironmentFile=/etc/lifeos/tts-server.env

# Phase 8: bind on the bridge interface so PublishPort= forwards reach the
# server. Without this the server binds only the container's own loopback,
# the bridge IP has nobody listening, and the daemon's loopback connections
# get silently refused.
Environment=LIFEOS_TTS_ENGINE_HOST=0.0.0.0

[Service]
# Defense-in-depth Capa 4: ensure image is present locally before starting.
ExecStartPre=/usr/local/bin/lifeos-ensure-images

# Capa 3 — auto-restart.
Restart=always
RestartSec=5s

# Kokoro warmup is now FAST because models are baked in (no model load
# round-trip to HF Hub, no big tensor copies between mounts).
TimeoutStartSec=60s

# Resource limits — applied to the systemd cgroup wrapping the podman
# process. With --cgroups=split the container gets its own sub-cgroup;
# these cap the whole service envelope.
MemoryMax=1536M
CPUQuota=200%

# DO NOT add `NoNewPrivileges=true`. NNP + Quadlet `--cgroups=split` on
# Fedora 44 + podman 5.x blocks /run/libpod/alive.lck acquisition.
# Container-side hardening goes in [Container] (Quadlet translates to
# podman flags directly).

[Install]
WantedBy=multi-user.target default.target
