# containers/lifeos-tts/lifeos-tts.container
# lifeos-tts.container — Quadlet for Kokoro TTS server (self-contained image)
#
# Phase 6 of the architecture pivot — the side image now carries the full
# Kokoro venv + voice models + voices-manifest.json + the server script.
# The bootc image no longer ships /opt/lifeos/kokoro-env or
# /usr/local/bin/lifeos-tts-server.py — both moved INTO the container
# image, removing ~850 MB from bootc.
#
# PRD: docs/strategy/prd-architecture-pivot-lean-bootc-quadlet.md

[Unit]
Description=LifeOS Kokoro TTS Server (containerized)
Documentation=https://github.com/hexgrad/kokoro
After=network-online.target lifeos-image-guardian.service
Wants=network-online.target
# After Phase 3 flips lifeosd to a Quadlet, this should re-target the
# Quadlet-generated unit name. For now legacy lifeosd.service is the
# host unit; ordering against a non-existent unit is silently dropped.
Before=lifeos-lifeosd.service

[Container]
Image=ghcr.io/hectormr206/lifeos-tts:stable
ContainerName=lifeos-tts
HostName=lifeos-tts

# Network=host — Phase 1 networking choice (PRD §3 Opción A). Container
# binds 127.0.0.1:8084 in the host netns; lifeosd reaches it unchanged.
Network=host

# Only surface left: the env file. Phase 6 makes everything else
# (venv, voice models, manifest, server script) live INSIDE the
# container image.
#
# We use Quadlet's [Container] EnvironmentFile= which translates to
# `podman --env-file=/etc/lifeos/tts-server.env`, propagating the
# file's vars into the container process env. The script reads via
# os.environ, never opens the path inside the container — so a
# Volume= bind of the same file would be redundant AND would trigger
# an SELinux `:z` relabel of /etc/lifeos/tts-server.env to
# container_file_t (host services that read the same file would
# then get an unexpected MCS label). Single source of truth: this
# EnvironmentFile= line.
EnvironmentFile=/etc/lifeos/tts-server.env

[Service]
# Defense-in-depth Capa 4: ensure image is present locally before starting.
ExecStartPre=/usr/local/bin/lifeos-ensure-images

# Capa 3 — auto-restart.
Restart=always
RestartSec=5s

# Kokoro warmup is now FAST because models are baked in (no model load
# round-trip to HF Hub, no big tensor copies between mounts).
TimeoutStartSec=60s

# Resource limits — applied to the systemd cgroup wrapping the podman
# process. With --cgroups=split the container gets its own sub-cgroup;
# these cap the whole service envelope.
MemoryMax=1536M
CPUQuota=200%

# DO NOT add `NoNewPrivileges=true`. NNP + Quadlet `--cgroups=split` on
# Fedora 44 + podman 5.x blocks /run/libpod/alive.lck acquisition.
# Container-side hardening goes in [Container] (Quadlet translates to
# podman flags directly).

[Install]
WantedBy=multi-user.target default.target
