# containers/lifeos-tts/lifeos-tts.container
# lifeos-tts.container — Quadlet definition for Kokoro TTS server
#
# Phase 1 of the architecture pivot — first system Quadlet.
# PRD: docs/strategy/prd-architecture-pivot-lean-bootc-quadlet.md
#
# NOT auto-installed yet. When Phase 1 actually flips, this file will be
# copied into image/files/etc/containers/systemd/lifeos-tts.container by the
# Containerfile, and the legacy /etc/systemd/system/lifeos-tts-server.service
# will be removed (or masked) in the same change.
#
# To deploy this manually for testing on a laptop already on 0.8.25+:
#   1. Stop the legacy service:
#        sudo systemctl stop lifeos-tts-server.service
#        sudo systemctl mask lifeos-tts-server.service
#   2. Drop this file at /etc/containers/systemd/lifeos-tts.container.
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl start lifeos-tts.service
#   5. curl http://127.0.0.1:8084/health
#
# Both the legacy service and this Quadlet bind 127.0.0.1:8084 — running
# them simultaneously will fail one. Always mask one before starting the other.

[Unit]
Description=LifeOS Kokoro TTS Server (containerized)
Documentation=https://github.com/hexgrad/kokoro
After=network-online.target lifeos-image-guardian.service
Wants=network-online.target
# lifeosd reads from :8084 — keep TTS up before lifeosd. Quadlets generate
# .service units, so reference the generated name.
Before=lifeosd.service

[Container]
Image=ghcr.io/hectormr206/lifeos-tts:stable
ContainerName=lifeos-tts
HostName=lifeos-tts

# Network=host is the Phase 1 choice (PRD §3 — Networking Opción A).
# - Container listens on 127.0.0.1:8084 inside its netns AND on the host's
#   loopback simultaneously, so existing clients (lifeosd) reach it without
#   changing any URL.
# - Network=host loses isolation between containers but they're already all
#   on localhost, so we lose nothing real.
# - Phase 6 will migrate to a podman bridge `lifeos-net` for proper isolation.
Network=host

# Bind-mount the live script from the host as read-only. This means the
# container image is stable across script edits — change .py, systemctl
# restart lifeos-tts, see new behavior. No rebuild, no push, no pull.
# `:Z` relabels for SELinux so the container can read the file.
Volume=/usr/local/bin/lifeos-tts-server.py:/usr/local/bin/lifeos-tts-server.py:ro,Z

# The legacy systemd unit reads /etc/lifeos/tts-server.env. Same here for
# parity — the env file lives in the bootc image, no container rebuild
# required when env keys change.
Volume=/etc/lifeos/tts-server.env:/etc/lifeos/tts-server.env:ro,Z

# Defense-in-depth Capa 4: re-pull the image from the registry if it was
# deleted from local storage. Cheap if image is present (skopeo no-op).
ExecStartPre=/usr/local/bin/lifeos-ensure-images

[Service]
# Capa 3 of defense-in-depth — auto-restart with watchdog.
Restart=always
RestartSec=5s

# Kokoro warmup loads ~500 MB of model weights — give it room to spin up.
TimeoutStartSec=60s

# WatchdogSec relies on the container reporting back to systemd via
# sd_notify. The current Python script doesn't do that, so we leave Watchdog
# off for now and rely on health checks (planned: HTTP probe on /health).
# When the script learns sd_notify, uncomment:
# WatchdogSec=60s

# Resource limits — same as the legacy systemd unit.
MemoryMax=1536M
CPUQuota=200%

# Security hardening that doesn't conflict with podman containerization.
# Several hardenings (ProtectSystem=strict, PrivateDevices=true) are
# implicit when running inside a container.
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target default.target
