# containers/lifeos-lifeosd/lifeos-lifeosd.container
# lifeos-lifeosd.container — Quadlet for the Axi daemon
#
# Phase 3 of the architecture pivot. NOT auto-installed yet.
# The "core" daemon — HTTP API, SimpleX/Telegram bridges, SQLite memory.
# The GTK4 desktop overlay is OUT OF SCOPE for this container (Phase 3b).

[Unit]
Description=LifeOS Axi Daemon (containerized core)
Documentation=https://github.com/hectormr206/lifeos
# All Quadlet-generated unit names use the `lifeos-` prefix. The earlier
# version of this file referenced `llama-server.service` and
# `simplex-chat.service` (the LEGACY non-prefixed host services); systemd
# silently dropped ordering against those after Phase 4/5 flipped the
# legacy units to masked/removed state.
After=network-online.target lifeos-image-guardian.service lifeos-state-migrate.service lifeos-net-network.service lifeos-tts.service lifeos-llama-embeddings.service lifeos-llama-server.service lifeos-simplex-bridge.service
# Hard bridge dependency — see lifeos-tts.container for rationale. Without
# this lifeos-net.network failing to start would let the daemon come up
# with no network attachment, PublishPort= forwards die silently, and the
# dashboard/sentinel/CLI all hang on connect.
Requires=lifeos-net-network.service
# `Wants=` on the migration unit (NOT `Requires=`) — round-2 audit caught
# that `Requires=` would block the daemon forever if the migration unit
# failed even once (Type=oneshot has no Restart= and the daemon would
# need a manual `systemctl reset-failed` to come back). With `Wants=` the
# migration is best-effort: if it fails, the daemon still starts against
# whatever state /var/lib/lifeos has, the operator gets a journal entry
# explaining the failure, and the sentinel can do its job. The daemon's
# memory-plane code already handles "no memory.db yet" cleanly (fresh
# DB on first start) — losing migration is recoverable; losing the
# entire daemon-start path is not.
Wants=lifeos-state-migrate.service network-online.target lifeos-tts.service lifeos-llama-embeddings.service lifeos-llama-server.service lifeos-simplex-bridge.service

[Container]
Image=ghcr.io/hectormr206/lifeos-lifeosd:stable
# Phase 6 — auto-update: podman-auto-update.timer (daily) compares the
# registry digest for the pinned tag and rolling-restarts the service if
# it changed.
AutoUpdate=registry
ContainerName=lifeos-lifeosd
HostName=lifeos-lifeosd

# Phase 8b: UDS+SO_PEERCRED migration. The daemon now exposes two listeners:
#
#   1. Unix-domain socket at /run/lifeos/lifeosd.sock (host-shared via Volume=
#      below). Machine clients (life CLI, sentinel.sh, update-*.sh) connect
#      here. The kernel's SO_PEERCRED gate at accept(2) time authenticates
#      callers by UID — no loopback-IP middleware needed.
#
#   2. TCP on 0.0.0.0:8081 inside the container (browser dashboard surface).
#      PublishPort=127.0.0.1:8081:8081 forwards this to the host loopback ONLY
#      — Firefox can reach http://127.0.0.1:8081/dashboard/ while sibling
#      containers on lifeos-net CANNOT reach port 8081 (PublishPort binds the
#      host interface, not the bridge). This closes the R9 lateral-move risk
#      and makes the old require_loopback_peer middleware obsolete (deleted in
#      daemon/src/api/mod.rs Phase 8b commit 3).
#
# lifeos-net.network gives the daemon DNS access to sibling containers
# (lifeos-tts, lifeos-llama-server, etc.) via their container names, which
# matches the URL defaults in endpoints.rs. No environment overrides needed.
Network=lifeos-net.network
PublishPort=127.0.0.1:8081:8081

# === Bind mounts: state lives on host, container is stateless ===
#
# /var/lib/lifeos — encrypted memory DB, session_store, calendar, task
# queue, scheduled tasks, Kokoro voice cache, etc. RW. C3 fix means
# container runs as root inside, so the host UID 1000 owner is reachable.
#
# ':z' (lowercase, SHARED MCS label) — multiple containers and the
# legacy host services bind paths under /var/lib/lifeos and /etc/lifeos.
# ':Z' (uppercase, private) would assign a per-container category and
# lock out everyone else — exactly the rollback-breaking behavior the
# audit flagged.
Volume=/var/lib/lifeos:/var/lib/lifeos:z

# (Models bind is REDUNDANT — /var/lib/lifeos/models is already covered
# by the parent /var/lib/lifeos:z mount above. A nested RO mount over
# an RW parent has unreliable semantics across podman versions and
# triggers double-relabel races. Removed.)

# Configs read-only. /etc/lifeos has tts-server.env, llama-server.env,
# llm-providers.env, etc. ':z' shared label so TTS/embeddings/legacy
# services keep access.
Volume=/etc/lifeos:/etc/lifeos:ro,z

# Telegram bridge in lifeosd needs CA bundle to verify TLS to
# api.telegram.org. ':z' shared label — touching the system-wide CA
# trust dir with ':Z' would relabel it and break TLS for every other
# host service.
Volume=/etc/pki/ca-trust/extracted:/etc/pki/ca-trust/extracted:ro,z

# === Critical environment variables ===
#
# Round-2 JD ordering note: systemd applies Environment= and EnvironmentFile=
# in DIRECTIVE ORDER (later overrides earlier). EnvironmentFile= MUST come
# FIRST so the inline Environment= values below win for any LIFEOS_* key
# the legacy llm-providers.env file might also contain — otherwise an
# operator who copied LIFEOS_TTS_URL into /etc/lifeos/llm-providers.env
# could silently override the bridged Quadlet topology.
#
# Quadlet does NOT honor the systemd `-` "optional" prefix on
# EnvironmentFile= when the directive lives in the [Container] section —
# it interprets the dash as a leading path segment and tries to open the
# file as `/usr/share/containers/systemd/-/etc/lifeos/...`. The host
# bootc image always ships /etc/lifeos/llm-providers.env (lifeos-secrets.conf
# enforces 0600 root) so the file is always present — drop the dash.
EnvironmentFile=/etc/lifeos/llm-providers.env
Environment=LIFEOS_DASHBOARD_DIR=/usr/share/lifeos/dashboard
Environment=LIFEOS_RUNTIME_DIR=/run/lifeos
Environment=RUST_LOG=info,lifeosd=debug
# Phase 8b: UDS API socket path and UID allowlist.
# LIFEOS_API_SOCKET: path to the Unix-domain socket for machine clients.
#   Must match the Volume= bind-mount below (/run/lifeos:/run/lifeos:z).
# LIFEOS_API_UID: UID of the host user allowed to connect via UDS in
#   addition to root (UID 0). Default 1000 = the lifeos system user.
#   Both variables use the same names as the daemon's env-var parsing.
Environment=LIFEOS_API_SOCKET=/run/lifeos/lifeosd.sock
Environment=LIFEOS_API_UID=1000
# LIFEOS_API_BIND is intentionally unset. The daemon's default
# (0.0.0.0:8081) is correct for the bridge topology: PublishPort=
# forwards to host loopback only, so the wildcard bind inside the
# container does not expose the port externally.

# Sibling container URLs default to http://127.0.0.1:80xx in endpoints.rs.
# Under Network=lifeos-net.network the siblings' PublishPort= entries forward
# their ports back to the host loopback, so these URLs resolve correctly.
# Container-name DNS (e.g. http://lifeos-tts:8084) is available as well.

# Memory DB encryption key derives from /etc/machine-id + /etc/hostname
# (see daemon/src/memory_plane.rs::derive_machine_key). Without these
# bind-mounts, podman generates a fresh /etc/machine-id inside the
# container, the derived key changes, and the existing memory.db (encrypted
# with the host's key) becomes unreadable. RO + no `:z` because relabeling
# /etc/machine-id to container_file_t breaks D-Bus + journald + half of
# systemd on the host.
Volume=/etc/machine-id:/etc/machine-id:ro
Volume=/etc/hostname:/etc/hostname:ro

# /run/lifeos shared with host: the daemon writes
# /run/lifeos/bootstrap.token (LIFEOS_RUNTIME_DIR=/run/lifeos), and
# lifeos-sentinel.sh on the host reads that file every 30 s to authenticate
# its health probe. Without this bind the daemon writes to the container's
# private /run tmpfs; the sentinel sees an empty token, fails every health
# check, and crash-loops the daemon. Host-side /run/lifeos is created at
# boot by tmpfiles.d/lifeos-runtime.conf. ':z' (shared MCS) — same path is
# also touched by other host-side LifeOS services in future phases.
Volume=/run/lifeos:/run/lifeos:z

# === GPU access for hardware-tuned profile detection ===
#
# The daemon's ai_runtime_profile.rs::compute_runtime_profile reads
# /sys/class/drm + nvidia-smi to detect a dedicated GPU and pick the
# right Qwen model (9B + GPU layers if dedicated GPU present, smaller
# 4B CPU profile if not). Pre-Phase-3 lifeosd ran on the host and saw
# /dev/nvidia0 directly. Post-Phase-3 it lives in this container with
# no device access — so it always reports "no dedicated GPU", writes
# /var/lib/lifeos/llama-server-runtime-profile.env without
# LIFEOS_AI_MODEL/LIFEOS_AI_MMPROJ, and lifeos-llama-server.service
# then starts with an empty `--model` flag and crash-loops.
#
# CDI passthrough fixes this — same mechanism lifeos-llama-server
# uses. The daemon doesn't actually do GPU compute (that's
# llama-server's job); it just needs to *detect* the GPU. The host
# nvidia-container-toolkit publishes nvidia.yaml at /run/cdi or
# /etc/cdi; podman injects /dev/nvidia* + libnvidia-* into the
# container at start. Skipped silently on hosts with no NVIDIA driver
# (CDI spec absent) — daemon falls back to CPU profile.
AddDevice=nvidia.com/gpu=all

# === Container hardening (Quadlet [Container] directives) ===
# Strategy: drop everything by default, add only what the daemon needs.
DropCapability=all
# DAC_OVERRIDE: read/write host-owned /var/lib/lifeos/*.db.
# CHOWN, FOWNER: lifeosd creates new files under /var/lib/lifeos and may
#   chown to match existing perms.
AddCapability=DAC_OVERRIDE CHOWN FOWNER
# SELinux type — container_t is podman default for rootful but explicit
# is documentation. Per-container MCS isolation: add SecurityLabelLevel=
# if/when sibling containers stop sharing :z paths.
SecurityLabelType=container_t

[Service]
# Defense-in-depth Capa 4: ensure image is present locally before starting.
# Quadlet requires ExecStartPre in [Service], not [Container].
ExecStartPre=/usr/local/bin/lifeos-ensure-images

# CDI spec poll — same guard the lifeos-llama-server Quadlet uses. The
# AddDevice=nvidia.com/gpu=all directive above only succeeds when CDI
# has published nvidia.yaml; on cold boot nvidia-cdi-refresh.service
# may not have run yet (we can't After= it without producing the
# ordering cycle that bit lifeos-llama-server in production). Up to
# 60 s, fall through if absent — Restart=always covers a hard miss
# and on hosts without NVIDIA the loop simply times out and the daemon
# proceeds with no GPU detected (CPU profile path).
ExecStartPre=-/bin/bash -c 'i=0; while [ $i -lt 60 ]; do for f in /run/cdi/nvidia.yaml /var/run/cdi/nvidia.yaml /etc/cdi/nvidia.yaml; do [ -s "$f" ] && exit 0; done; sleep 1; i=$((i+1)); done; exit 0'

# Capa 3 of defense-in-depth — auto-restart with watchdog.
Restart=always
RestartSec=5s

# Daemon initial startup runs migrations + opens encrypted memory DB +
# bootstraps SimpleX bridge — give it time on cold boot.
TimeoutStartSec=120s

# WatchdogSec is INTENTIONALLY NOT SET.
# Audit (Día del Juicio) found: daemon/src/main.rs implements notify_ready()
# (one-shot READY=1) and notify_stopping() (one-shot STOPPING=1) but NO
# periodic loop sending WATCHDOG=1. Setting WatchdogSec=300s here would
# make systemd SIGABRT the container every 5 minutes waiting for a
# heartbeat that never comes — RestartSec=5s would re-launch and the
# cycle repeats forever.
#
# Re-enable ONLY after daemon/src/main.rs implements the watchdog ping
# loop following the systemd protocol (every WatchdogSec/2 seconds):
# WatchdogSec=300s

# Resource limits — applied to the systemd cgroup wrapping the podman
# process. With Quadlet's --cgroups=split the container gets its own
# sub-cgroup; this caps the whole service envelope.
MemoryMax=1G

# DO NOT add NoNewPrivileges= here. NNP + Quadlet --cgroups=split on
# Fedora 44 + podman 5.x blocks /run/libpod/alive.lck acquisition
# (broke TTS in PR #77 with the same daemon stack). Container hardening
# lives in [Container] above (DropCapability/AddCapability/SecurityLabel).

[Install]
WantedBy=multi-user.target default.target
