# image/files/etc/systemd/system/llama-embeddings.service.d/20-lifeos-background.conf
# Pin llama-embeddings to the LifeOS background slice (E-cores, low priority).
# Embedding generation runs continuously but is never urgent — E-cores
# handle it efficiently without waking P-cores.
[Service]
Slice=lifeos-background.slice
Nice=15
IOSchedulingClass=idle
