# image/files/etc/systemd/system/simplex-chat.service.d/10-bootstrap-profile.conf
# LifeOS dev-mode override: bootstrap the SimpleX profile inline.
#
# The in-image `lifeos-simplex-setup.sh` ships in /usr/local/bin which
# is read-only on bootc, so the dev-mode loop cannot iterate on that
# script without a full image rebuild. This dropin injects the same
# "create the bot profile on first run by piping a display name"
# logic as an additional ExecStartPre, so restart-looped simplex-chat
# units can be fixed by sudo install + daemon-reload + restart.
#
# Once the next image rebuild bakes the fix into /usr/local/bin/
# lifeos-simplex-setup.sh this dropin becomes a redundant no-op but
# stays harmless because the setup marker (~/simplex/.setup-done)
# short-circuits it.
[Service]
ExecStartPre=/bin/bash -c '\
    set -eu; \
    D="/var/lib/lifeos/simplex"; \
    P="$$D/bot"; \
    M="$$D/.bootstrap-dropin-done"; \
    if [ -f "$$M" ] && [ -f "$${P}_chat.db" ]; then exit 0; fi; \
    mkdir -p "$$D"; \
    chmod 0700 "$$D"; \
    rm -f "$${P}_chat.db" "$${P}_agent.db"; \
    printf "Axi\\n" | /usr/local/bin/simplex-chat -d "$$P" -y -e "/quit" >/dev/null 2>&1 || true; \
    if [ ! -f "$${P}_chat.db" ]; then \
        echo "simplex bootstrap: database never created at $${P}_chat.db" >&2; \
        exit 1; \
    fi; \
    touch "$$M"; \
'
