# image/files/etc/systemd/system/llama-server.service.d/95-user-override.conf
# User-managed override file for llama-server env vars.
#
# Loaded AFTER:
#   - /etc/lifeos/llama-server.env                           (baseline)
#   - /var/lib/lifeos/llama-server-runtime-profile.env       (hardware-tuned)
#   - /var/lib/lifeos/llama-server-game-guard.env            (game guard)
# so any LIFEOS_AI_* value here ALWAYS wins.
#
# The leading dash makes it a no-op when the file is missing — the user
# only writes it via the dashboard ("LLM Configuration" → "Save & Restart")
# or via POST /api/v1/llm/ctx-size.
#
# Mode 0644, owner lifeos:lifeos. Safe to delete by hand to revert to the
# hardware-tuned profile.
[Service]
EnvironmentFile=-/var/lib/lifeos/llama-server-user-override.env
